# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=python -m uvicorn backend.api:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/fastapi.log
stderr_logfile=/tmp/fastapi.log
environment=PYTHONUNBUFFERED="1"

[program:streamlit]
command=python -m streamlit run streamlit_app.py
    --server.port=7860
    --server.address=0.0.0.0
    --server.headless=true
    --server.enableCORS=false
    --server.enableXsrfProtection=false
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/streamlit.log
stderr_logfile=/tmp/streamlit.log
environment=PYTHONUNBUFFERED="1",BACKEND_URL="http://localhost:8000"
